# wifi-provisioning/systemd/button-to-ap.service
[Unit]
Description=WiFi-hotspot: äänikortin nappi (GPIO27) käynnistää AP-tilan
After=multi-user.target
# Ei vaadi verkkoa; kuuntelee vain nappia

[Service]
Type=simple
WorkingDirectory=/opt/wifi-provisioning
ExecStart=/usr/bin/python3 /opt/wifi-provisioning/button-to-ap.py
Restart=always
RestartSec=5
User=root
Group=root

[Install]
WantedBy=multi-user.target
