# wifi-provisioning/systemd/wifi-watchdog.timer
[Unit]
Description=WiFi-hotspot: tarkista verkko 90 s v채lein, k채ynnist채 AP jos ei yhteytt채
After=network-online.target

[Timer]
OnBootSec=1min
OnUnitActiveSec=90
Persistent=yes

[Install]
WantedBy=timers.target
